# dyndns.service
[Unit]
Description=DynDNS Updater

[Service]
Type=oneshot
ExecStart=/etc/dyndns/dyndns-update
PrivateTmp=true

